# docker/rootfs/etc/services.d/tc66c-mqtt/finish
#!/bin/sh
# shellcheck shell=sh

s6-svscanctl -t /var/run/s6/services 2> /dev/null
