# docker/rootfs/etc/services.d/tc66c-mqtt/finish
#!/bin/ash
# shellcheck shell=ash

set -euo pipefail

s6-svscanctl -t /var/run/s6/services 2> /dev/null
